# deploy/systemd/pk-reminders.timer
[Unit]
Description=Daily Personal Knowledge reminder scan

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true

[Install]
WantedBy=timers.target
